# launchd/com.stockcharts.daily-analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockcharts.daily-analysis</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/charlie/Documents/projects/stockcharts/.venv/bin/python</string>
        <string>-m</string>
        <string>src.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/charlie/Documents/projects/stockcharts</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/charlie/Documents/projects/stockcharts/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/charlie/Documents/projects/stockcharts/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Run job if Mac was asleep at scheduled time -->
    <key>StartInterval</key>
    <integer>0</integer>
</dict>
</plist>
